# Docker/crontab
*/10 * * * * cd /app && /usr/local/bin/flask update_infogen
